# deploy/google-cloud/systemd/caio-sqlite-backup.timer
[Unit]
Description=Back up Caio SQLite database hourly

[Timer]
OnCalendar=hourly
Persistent=true
Unit=caio-sqlite-backup.service

[Install]
WantedBy=timers.target
